# .local/var/run/runit/pipewire-pulse/run
#!/bin/sh
sv check ../wireplumber >/dev/null || exit 1
exec pipewire-pulse
