# .local/var/run/runit/pipewire/finish
#!/bin/sh
sv d ../pipewire-pulse ../wireplumber
