# .local/var/run/runit/pipewire/run
#!/bin/sh
sv u ../pipewire-pulse ../wireplumber
exec pipewire
