# .local/var/run/runit/pipewire/run
#!/bin/sh
sv u ../pipewire-pulse ../wireplumber 2>/dev/null || exit
exec pipewire
