# .local/var/run/runit/wireplumber/run
#!/bin/sh
sv check ../pipewire >/dev/null || exit 1
exec wireplumber
