# .local/var/run/runit/libretranslate/run
#!/bin/sh
mkdir -p ~/.cache/libretranslate
cd ~/.cache/libretranslate || exit 1
exec libretranslate --threads "$(nproc)" >/dev/null 2>&1
