# .local/var/run/runit/syncthing/run
#!/bin/sh
sv check ../netCheck >/dev/null || exit 1
exec 2>&1
exec syncthing --logflags 0 --no-browser
