# .local/var/run/runit/deluged/run
#!/bin/sh
sv check ../netCheck >/dev/null || exit 1

exec bwrap.sh -more -net -exec --bind 'printf "$HOME/%s\n" download vids/tv music .config/deluge' deluged -d
