# .local/var/run/runit/deluged/run
#!/bin/sh
sv check ../netCheck >/dev/null || exit 1
exec deluged -d
