# .local/var/run/runit/deluged/run
#!/bin/sh
sv check ../netCheck >/dev/null || exit 1

exec bwrap.sh -noshare -net -exec --bind 'printf "$HOME/%s\0" download vids/tv music .config/deluge' -- deluged -d
