# .local/var/run/runit/deluged/run
#!/bin/sh
sv check ../netCheck || exit 1
exec deluged -d
