# .local/var/run/runit/mpd/run
#!/bin/sh
timeout 2 sh -c 'until mpc pause; do true; done' &
exec mpd --no-daemon
