# .local/var/run/runit/firefox-ram/finish
#!/bin/sh
. ./conf
memory-sync firefox "$directory"
