# .local/var/run/runit/firefox-ram/run
#!/bin/sh
exec memory-sync -n -D -u 30m -d 0 -t 5 firefox "$HOME/.mozilla/firefox" firefox
