# .local/var/run/runit/firefox-ram/run
#!/bin/sh
. ./conf
exec memory-sync -D -u 30m -d 0 -t 5 firefox "$directory" "$binary" >/dev/null 2>&1
