# .local/var/run/runit/xapps/finish
#!/bin/sh
. ./conf
for serv in $services; do
	sv d "../$serv"
done
